# tor-migration-package/systemd/tor@default.service.d/proxy.conf
[Service]
EnvironmentFile=-/etc/tor/proxy.env
